# APP/etc/systemd/system/tale.service

[Unit]
Description=tale
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tale.sh
Restart=always

[Install]
WantedBy=multi-user.target
